# etc/supervisor/conf.d/shortie-dev.conf
[supervisord]

[program:shortie-api-dev]
command = gunicorn  shortie.api:app -b localhost:9001
